# azure/61-cloud-netconfig-hotplug.rules
SUBSYSTEM=="net", KERNEL=="eth*", DRIVERS=="hv_netvsc", RUN+="/usr/libexec/netconfig/scripts/cloud-netconfig-hotplug"
