# systemd/cloud-netconfig.timer
[Unit]
Description=Periodically run cloud-netconfig

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
