# systemd/cloud-netconfig.service
[Unit]
Description=Update cloud-netconfig Setup

[Service]
Type=oneshot
ExecStart=-/bin/bash -c '/sbin/netconfig -m cloud-netconfig update'
